# roles/daemons/templates/supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:etcd]
command=/home/{{ ansible_user }}/tmp/etcd-v3.2.2-linux-amd64/etcd
autostart=true
autorestart=false
stderr_logfile=/var/log/etcd.err.log
stdout_logfile=/var/log/etcd.out.log

#[program:dolphin]
#command={{ cis_root }}/dolphin/dolphin
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/dolphin.err.log
#stdout_logfile=/var/log/dolphin.out.log
#
#[program:jellyfish]
#command={{ cis_root }}/jellyfish/jellyfish
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/jellyfish.err.log
#stdout_logfile=/var/log/jellyfish.out.log
#
#
#[program:whale]
#command={{ cis_root }}/whale/whale
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/whale.err.log
#stdout_logfile=/var/log/whale.out.log
#
#[program:poseidon]
#command={{ cis_root }}/poseidon/poseidon
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/poseidon.err.log
#stdout_logfile=/var/log/poseidon.out.log
#
#[program:seabed]
#command={{ cis_root }}/seabed/seabed
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/seabed.err.log
#stdout_logfile=/var/log/seabed.out.log
#
#[program:seahorse]
#command={{ cis_root }}/seahorse/seahorse
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/seahorse.err.log
#stdout_logfile=/var/log/seahorse.out.log
#
#[program:mackerel]
#command={{ cis_root }}/mackerel/mackerel
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/mackerel.err.log
#stdout_logfile=/var/log/mackerel.out.log
#
#[program:viettel]
#command={{ cis_root }}/viettel/viettel
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/viettel.err.log
#stdout_logfile=/var/log/viettel.out.log
#
#[program:vietnamobile]
#command={{ cis_root }}/vietnamobile/vietnamobile
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/vietnamobile.err.log
#stdout_logfile=/var/log/vietnamobile.out.log
#
#[program:octopus]
#command={{ cis_root }}/octopus/octopus
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/octopus.err.log
#stdout_logfile=/var/log/octopus.out.log
#
#[program:seaweed]
#command={{ cis_root }}/seaweed/seaweed
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/seaweed.err.log
#stdout_logfile=/var/log/seaweed.out.log

[program:seaweedfs_master]
command=/home/{{ ansible_user }}/seaweedfs/weed master -defaultReplication='000' -ip='localhost' -ip.bind='0.0.0.0' -port='8898' -maxCpu=2 -mdir='/home/{{ ansible_user }}/seaweedfs/master/metadata' -peers='localhost:7898' -volumePreallocate -volumeSizeLimitMB='1024' -whiteList='localhost,127.0.0.1,::1'
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/seaweedfs_master.err.log
stdout_logfile=/var/log/seaweedfs_master.out.log

[program:seaweedfs_vol1]
command=/home/{{ ansible_user }}/seaweedfs/weed volume -dataCenter='dc1' -rack='rack1' -ip='localhost' -ip.bind='0.0.0.0' -port='8899' -dir='/home/{{ ansible_user }}/seaweedfs/volume1/data' -images.fix.orientation='false' -index='leveldb' -max=20 -maxCpu=2 -mserver='localhost:8898' -whiteList='localhost,127.0.0.1,::1'
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/seaweedfs_vol1.err.log
stdout_logfile=/var/log/seaweedfs_vol1.out.log

[program:seaweedfs_vol2]
command=/home/{{ ansible_user }}/seaweedfs/weed volume -dataCenter='dc1' -rack='rack1' -ip='localhost' -ip.bind='0.0.0.0' -port='9988' -dir='/home/{{ ansible_user }}/seaweedfs/volume1/data' -images.fix.orientation='false' -index='leveldb' -max=20 -maxCpu=2 -mserver='localhost:8898' -whiteList='localhost,127.0.0.1,::1'
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/seaweedfs_vol2.err.log
stdout_logfile=/var/log/seaweedfs_vol2.out.log

#[program:napas]
#command={{ cis_root }}/napas/napas
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/napas.err.log
#stdout_logfile=/var/log/napas.out.log


[program:smssender]
command={{ cis_root }}/seapostman/SmsSender/SmsSender
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/smssender.err.log
stdout_logfile=/var/log/smssender.out.log

[program:emailsender]
command={{ cis_root }}/seapostman/EmailSender/EmailSender
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/emailsender.err.log
stdout_logfile=/var/log/emailsender.out.log


#[program:mobifone]
#command={{ cis_root }}/mobifone/mobifone
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/mobifone.err.log
#stdout_logfile=/var/log/mobifone.out.log


#[program:napas_fasthttp]
#command={{ cis_root }}/napas_fasthttp/napas_fasthttp
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/napas_fasthttp.err.log
#stdout_logfile=/var/log/napas_fasthttp.out.log
#
#[program:xevipham_fasthttp]
#command={{ cis_root }}/xevipham_fasthttp/xevipham_fasthttp
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/xevipham_fasthttp.err.log
#stdout_logfile=/var/log/xevipham_fasthttp.out.log
#
#[program:vinaphone_fasthttp]
#command={{ cis_root }}/vinaphone_fasthttp/vinaphone_fasthttp
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/vinaphone_fasthttp.err.log
#stdout_logfile=/var/log/vinaphone_fasthttp.out.log
#
#[program:mobifone_fasthttp]
#command={{ cis_root }}/mobifone_fasthttp/mobifone_fasthttp
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/mobifone_fasthttp.err.log
#stdout_logfile=/var/log/mobifone_fasthttp.out.log
#
#[program:vietnamobile_fasthttp]
#command={{ cis_root }}/vietnamobile_fasthttp/vietnamobile_fasthttp
#user={{ ansible_user }}
#autostart=true
#autorestart=false
#stderr_logfile=/var/log/vietnamobile_fasthttp.err.log
#stdout_logfile=/var/log/vietnamobile_fasthttp.out.log
#
#[program:viettel_fasthttp]
#command={{ cis_root }}/viettel_fasthttp/viettel_fasthttp
user={{ ansible_user }}
autostart=true
autorestart=false
stderr_logfile=/var/log/viettel_fasthttp.err.log
stdout_logfile=/var/log/viettel_fasthttp.out.log

